# templates/etc/sv/home/run
#!/bin/bash
exec 2>&1
exec chpst -u<%= @user %> runsvdir <%= @home %>/service/enabled
